# systemd/fivem-ddos-failure@.service
[Unit]
Description=FiveM DDoS Monitor Failure Alert

[Service]
Type=oneshot
ExecStart=/usr/local/bin/fivem-ddos-failure-alert.sh
