# systemd/fivem-ddos.service
[Unit]
Description=FiveM DDoS Detection Service
